# templates/victoriametrics.service.j2
[Unit]
Description=VictoriaMetrics time-series database
After=network.target
RequiresMountsFor=/mnt/ExtData

[Service]
Type=simple
User=victoriametrics
Group=victoriametrics
ExecStart=/usr/local/bin/victoria-metrics-prod \
  -storageDataPath=/mnt/ExtData/victoriametrics \
  -httpListenAddr=127.0.0.1:8428 \
  -retentionPeriod=10y \
  -memory.allowedPercent=25
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
